# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log

[program:php5-fpm]
command=/usr/bin/php-fpm5 --nodaemonize

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
